# debian/wb-mqtt-gpio.service
[Unit]
Description=MQTT Driver for GPIO-controlled switches
Wants=wb-hwconf-manager.service wb-modules.service
After=network.target wb-hwconf-manager.service wb-modules.service

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/wb-mqtt-gpio
ExecStartPre=/usr/lib/wb-mqtt-gpio/generate-system-config.sh

[Install]
WantedBy=multi-user.target
Alias=wb-homa-gpio.service